# services/openinsider/openinsider-fetch-process/deploy/openinsider-cluster.timer
[Unit]
Description=Run OpenInsider cluster buys scraper every 6 hours

[Timer]
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
