# meta/systemd/openinsider-details.service
[Unit]
Description=OpenInsider insider details scraper (Phase 2)
After=network-online.target openinsider-cluster.service
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
Group=scraper
WorkingDirectory=/opt/valhalla/services/openinsider/openinsider-fetch-process
ExecStart=/opt/valhalla/services/openinsider/openinsider-fetch-process/venv/bin/python run_insider_details.py --all --limit 20
EnvironmentFile=/opt/valhalla/services/openinsider/openinsider-fetch-process/.env
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/valhalla/services/openinsider/openinsider-fetch-process/data /opt/valhalla/services/openinsider/openinsider-fetch-process/logs
PrivateTmp=true
