# meta/systemd/substack-pipeline.service
[Unit]
Description=Substack Newsletter Processing Pipeline
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/valhalla/services/substack/substack-fetch-process
ExecStart=/opt/valhalla/services/substack/substack-fetch-process/venv/bin/python orchestrate.py --stage all
EnvironmentFile=/opt/valhalla/services/substack/substack-fetch-process/.env

# Graceful shutdown -- sends SIGTERM, waits 30s before SIGKILL
TimeoutStopSec=30
KillSignal=SIGTERM

# Logging (stdout/stderr captured by journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=substack-pipeline

[Install]
WantedBy=multi-user.target
